# images/nginx/nginx-supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:php-fpm]
command = /usr/bin/php-fpm
user = root
autostart = true