# images/mariadb/mariadb-supervisor.ini
[supervisord]
nodaemon = true

[program:mariadb]
command=sh /run.sh
numprocs=1
autostart=true
autorestart=true 